# config/launchd/com.jun.tmux-dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jun.tmux-dev</string>
  <!-- At login, make sure a detached tmux "dev" session exists to attach to.
       Idempotent: only creates it if it isn't already running. Starting the
       tmux server here also lets tmux-continuum restore the last saved layout. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>/opt/homebrew/bin/tmux has-session -t dev 2>/dev/null || /opt/homebrew/bin/tmux new-session -d -s dev</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
